# .config/systemd/user/tfresync.timer
[Unit]
Description=Timer for tfsync 1Hour
PartOf=tfresync.service tfsync.service

[Timer]
OnUnitActiveSec=3h
